# src/wol-monitor.service.in
[Unit]
Description=Wake-on-LAN monitor service
After=network.target network-online.target

[Service]
Type=simple
Environment=WOL_GROUP=root
EnvironmentFile=@SYSCONFDIR@/default/wol-monitor
ExecStart=@LIBEXECDIR@/wol-monitor "$WOL_GROUP"
PrivateTmp=true

[Install]
WantedBy=multi-user.target
